# templates/default/poirot_receiver.conf.erb
start on runlevel [2345]
stop on runlevel [016]
respawn

exec /opt/poirot_receiver/current/bin/poirot_receiver \
  -port=<%= node['poirot']['port'] %> \
  -index_prefix=<%= node['poirot']['index_prefix'] %> \
  -debug=<%= node['poirot']['debug'] %>
