# templates/default/poirot.service.erb
[Unit]
Description=Poirot receiver
After=elasticsearch.service

[Service]
Type=simple
User=<%= @user %>
WorkingDirectory=<%= @app_dir %>/current
ExecStart=<%= @app_dir %>/poirot.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

