# templates/redhat/poirot_receiver.conf.erb
start on runlevel [2345]
stop on runlevel [016]
respawn

pre-start script

bash << "EOF"
  mkdir -p <%= @logdir %>
  touch <%= @logdir %>/poirot_receiver.log
  chown -R <%= @user %> <%= @logdir %>
EOF

end script

exec su <%= @user %> --session-command='/opt/poirot_receiver/current/bin/poirot_receiver -port=<%= node['poirot']['receiver']['port'] %> -index_prefix=<%= node['poirot']['receiver']['index_prefix'] %> -debug=<%= node['poirot']['receiver']['debug'] %> >> <%= @logdir %>/poirot_receiver.log 2>&1'
