# templates/redhat/poirot-notifications.conf.erb
start on runlevel [2345]
stop on runlevel [016]
respawn
respawn limit 10 600
chdir <%= @app_dir %>/current
env HOME=<%= @app_dir %>/current
limit nofile 4096 4096

exec su <%= @user %> --session-command='export RAILS_ENV=production; export NEWRELIC_ENABLE=true; <%= @bundle_command %> exec bin/notification_service >> <%= @app_dir %>/shared/log/notifications.log 2>&1'

